# stravaapp.service
[Unit]
Description=Strava Web App Service
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/Amanda-Jeremaiah-William-Tori

# Environment variables
Environment="PATH=/home/ec2-user/Amanda-Jeremaiah-William-Tori/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Executable command
ExecStart=/home/ec2-user/Amanda-Jeremaiah-William-Tori/.venv/bin/gunicorn \
    --bind 0.0.0.0:8000 \
    --workers 4 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    --log-level info \
    app:app

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stravaapp

[Install]
WantedBy=multi-user.target